# laboratorios/laboratorio1/laboratorio1-parte4/supervisor_services.conf
[program:python]
command=python3 /src/main.py
autostart=true
autorestart=true

# Logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
